# root/etc/s6-overlay/s6-rc.d/svc-output-log/run
#!/usr/bin/with-contenv bash
FIFO="/tmp/output-log.pipe"
rm -f "$FIFO"
mkfifo "$FIFO"
chown abc:abc "$FIFO"
exec cat <> "$FIFO"
